# deploy/systemd/opi5-safetyd.service
[Unit]
Description=OPi5 Edge Safety Daemon (local safety loop)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/edge-ai-safety-monitor.env
ExecStart=/opt/edge-ai-safety-monitor/build/opi5-controller/opi5_safetyd --mode loop
Restart=always
RestartSec=3
WorkingDirectory=/opt/edge-ai-safety-monitor
WatchdogSec=120

[Install]
WantedBy=multi-user.target
